# nixoraos-core/usr/lib/systemd/system/nixora-gpu-health.service
[Unit]
Description=Nixora GPU Health Check
After=graphical.target
Wants=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/lib/nixora/gpu-health-check.sh
TimeoutStartSec=30

[Install]
WantedBy=graphical.target